# docker/systemd/fair-edge.service
[Unit]
Description=Fair-Edge Sports Betting Analysis Platform
Documentation=https://github.com/your-org/fair-edge
Requires=docker.service
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=forking
RemainAfterExit=yes

# Working directory
WorkingDirectory=/opt/fair-edge

# Environment file
EnvironmentFile=/opt/fair-edge/.env.production

# Pre-start checks
ExecStartPre=/usr/bin/docker compose -f docker-compose.hetzner.yml config --quiet
ExecStartPre=/usr/bin/docker compose -f docker-compose.hetzner.yml pull --quiet

# Start the application stack
ExecStart=/usr/bin/docker compose -f docker-compose.hetzner.yml up -d

# Stop the application stack
ExecStop=/usr/bin/docker compose -f docker-compose.hetzner.yml down

# Reload configuration
ExecReload=/bin/bash -c "/usr/bin/docker compose -f docker-compose.hetzner.yml down && /usr/bin/docker compose -f docker-compose.hetzner.yml up -d"

# Resource limits
MemoryLimit=4G
CPUQuota=300%

# Security settings
User=root
Group=docker

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fair-edge

# Kill settings
TimeoutStartSec=300
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target